# cloudflow/workflows/scripts/system/munge.service.in
[Unit]
Description=MUNGE authentication service
Documentation=man:munged(8)
After=network.target
After=syslog.target
After=time-sync.target

[Service]
Type=forking
ExecStart=@sbindir@/munged --socket /run/munge/munge.socket.2 --key-file=/etc/munge/munge.key
PIDFile=/run/munge/munged.pid
User=munge
Group=munge
Restart=on-abort

[Install]
WantedBy=multi-user.target
